# umdloop_gui_web/deploy/systemd/cc/umdloop-gui.service
[Unit]
Description=UMD Loop GUI (Next.js on port 3000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/umdloop/infrastructure/umdloop_gui/umdloop_gui_web
ExecStart=/usr/bin/npx next start -p 3000
Restart=on-failure
RestartSec=2
Environment=NODE_ENV=production
User=umdloop

[Install]
WantedBy=multi-user.target
