# umdloop_gui_web/deploy/systemd/cc/umdloop-mission-sync.service
[Unit]
Description=UMD Loop Mission Sync Service (WebSocket on port 3001)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/umdloop/infrastructure/umdloop_gui/umdloop_gui_web
ExecStart=/usr/bin/python3 -m backend.mission_sync
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1
User=umdloop

[Install]
WantedBy=multi-user.target
