# umdloop_gui_web/deploy/systemd/orin/umdloop-backend.service
[Unit]
Description=UMD Loop Rover Backend (FastAPI on port 5000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/umdloop/umdloop_gui_web
ExecStart=/usr/bin/python3 -m umdloop_gui_web.backend
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
